# provisioning/roles/queue/templates/supervisor.ini.j2
[program:flower]
command = {{ anaconda_home }}/bin/celery flower 
	--broker=redis://guest:guest@localhost:6379/0

autostart=true
autorestart=true